# .root-fs/crontab
* * * * * timeout --signal=TERM 50s php82 /app/artisan schedule:run
